# config/bossa.service
[Unit]
Description=BOSSA edge runtime
After=network.target

[Service]
Type=notify
ExecStart=/opt/bossa/bin/bossa --config /etc/bossa/config.yaml
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
